# files/etc/cron.d/iptest
# ==========================================
# OpenWrt IPThrottle 插件 - Cron 定时任务
# 文件: /etc/cron.d/iptest
# 功能: 每分钟检查时间计划，自动启用/禁用规则
# 创建时间: 2026-06-09
# ==========================================

# 每分钟执行一次时间计划检查
* * * * * root /usr/sbin/iptest schedule

# 注意: schedule 命令会检查所有规则的 schedule_json，
# 根据当前时间和星期判断是否应该生效，必要时自动 reload
