# files/etc/uci-defaults/50-ipthrottle
#!/bin/sh
# IPThrottle UCI默认配置初始化
# 首次安装时执行

# 创建默认配置（如果不存在）
if [ ! -f /etc/config/ipthrottle ]; then
    cat > /etc/config/ipthrottle << 'EOF'
config ipthrottle
    option enabled '0'

config rule 'example1'
    option enabled '1'
    option name '示例规则-单IP限速'
    option wan_mask 'all'
    list ip_entry '192.168.1.100'
    option proto 'any'
    option mode 'independent'
    option upload_kbps '512'
    option download_kbps '2048'
    option priority '10'
    option schedule_type 'always'

config rule 'example2'
    option enabled '0'
    option name '示例规则-IP段限速'
    option wan_mask 'all'
    list ip_entry '192.168.1.10-192.168.1.20'
    option proto 'tcp'
    option mode 'shared'
    option upload_kbps '1024'
    option download_kbps '4096'
    option priority '20'
    option schedule_type 'weekly'
    list schedule_days '1'
    list schedule_days '2'
    list schedule_days '3'
    list schedule_days '4'
    list schedule_days '5'
    option schedule_start '09:00'
    option schedule_end '18:00'
    option comment '工作日白天限速'

config rule 'example3'
    option enabled '0'
    option name '示例规则-多IP独立限速'
    option wan_mask 'wan1'
    list ip_entry '192.168.1.100'
    list ip_entry '192.168.1.101-192.168.1.105'
    option proto 'any'
    option mode 'independent'
    option upload_kbps '256'
    option download_kbps '1024'
    option priority '30'
    option schedule_type 'always'
EOF
    echo "已创建默认配置文件 /etc/config/ipthrottle"
fi

# 启用服务（可选，取消注释以启用）
# /etc/init.d/ipthrottle enable
# /etc/init.d/ipthrottle start

exit 0
